# cron.d/v-debian-bug-467231-dyncomp
@daily root bash /etc/bin/cron/debian-bug-467231-dyncomp >/dev/null
